# scripts/agent-bot/mara.timer
[Unit]
Description=mara — poll prod heron every 5 minutes
After=network-online.target

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
# Survive a missed window (host asleep / restarted) by firing on resume.
Persistent=true

[Install]
WantedBy=timers.target
